# tf_file/files/etc/crontabs/root
58 * * * * logrotate /etc/logrotate.conf